# rootfs/etc/services.d/cron/run
#!/bin/execlineb -P
mkdir -p /etc/cron.d
exec cron -f -l -L 4
